# files/etc/services.d/10_dnsmasq/run
#!/usr/bin/execlineb -P
/usr/sbin/dnsmasq -d
