# files/etc/services.d/10_dnsmasq/run
#!/usr/bin/with-contenv sh
/usr/sbin/dnsmasq -d
