# files/etc/services.d/02_docker-gen/run
#!/usr/bin/execlineb -P
/usr/local/bin/docker-gen -watch -notify-output -notify "pkill dnsmasq" /etc/dnsmasq.tmpl /etc/dnsmasq.conf
